# templates/range-repair.service.erb
# Systemd unit file for Cassandra Range Repair
# Managed by Puppet, profile_ggonda_cassandr module

[Unit]
Description=Cassandra Range Repair Service
After=network.target cassandra.service
Requires=cassandra.service

[Service]
Type=forking
User=cassandra
Group=cassandra
WorkingDirectory=/var/lib/cassandra
PIDFile=/var/run/cassandra_range_repair.pid
ExecStart=/usr/local/bin/range-repair.sh start
ExecStop=/usr/local/bin/range-repair.sh stop
ExecReload=/usr/local/bin/range-repair.sh restart
TimeoutSec=600
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
